# root/etc/init.d/geph5
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1

CONF="geph5"
YAML_FILE="/var/etc/geph5.yaml"

start_service() {
    config_load $CONF
    local enabled
    config_get_bool enabled geph5 enabled 0

    # 【修复1】如果禁用，先强制清理旧进程，防止“僵尸”后台
    if [ "$enabled" -eq 0 ]; then
        killall -9 geph5-client >/dev/null 2>&1
        return 0
    fi

    local prog auth_token socks5_port http_port
    config_get prog geph5 binary_path "/usr/bin/geph5-client"
    config_get auth_token geph5 auth_token ""
    config_get socks5_port geph5 socks5_port "9909"
    config_get http_port geph5 http_port "9919"

    if [ ! -x "$prog" ]; then
        logger -t geph5 "Error: Binary not found at $prog"
        return 1
    fi

    if [ -z "$auth_token" ]; then
        logger -t geph5 "Error: Auth token is empty!"
        return 1
    fi

    # 【修复2】密钥防截断逻辑：将超长密钥切分为4段拼接
    local k1="3082010a0282010100d0ae53a794ea37bf2e100cb3a872177ec6c11e8375fdcbf92960ce0293465674eb1426a1841b7622a58979a5ff3f8aa2301a621545e9b90bb39d"
    local k2="1a6bfda19d6ca1aae74a3192ddfd2b9558eb652c3c2c22f42bdde272852fb67d93cae5846213512c474bf799844aee019bf718f6fa64223be06364459fc8dec66796b141d"
    local k3="450d730c4fffe1cac7df8f05591560afa44bcf274f6c0e2303b39c21ab09d19b459ee594512b8341f3d407c026e2509f42c6d89f82f6a3a36fd5c05ad423cd99ad390894"
    local k4="03eb9122ea60ef6648afff65438e8e26ce41fa55b9b18741965c77a627bae947bd38fc345e9adab42d6c458f6e194e4232cfd3f04924d5a5e932fe769610203010001"
    local bw_key="${k1}${k2}${k3}${k4}"

    mkdir -p /var/etc

    # 生成配置文件 (注入变量)
    cat <<YML > $YAML_FILE
socks5_listen: 127.0.0.1:$socks5_port
http_proxy_listen: 127.0.0.1:$http_port
exit_constraint: auto
bridge_mode: ForceBridges
cache: /tmp/geph5.db
passthrough_china: false
spoof_dns: false
task_limit: null
vpn: false
vpn_fd: null

credentials:
  secret: "$auth_token"

sess_metadata:
  filter:
    ads: true
    nsfw: false

broker:
  priority_race:
    '0':
      fronted:
        front: https://www.cdn77.com/
        host: 1826209743.rsc.cdn77.org
        override_dns: null
    '500':
      fronted:
        front: https://www.vuejs.org/
        host: svitania-naidallszei-3.netlify.app
        override_dns:
        - 75.2.60.5:443

broker_keys:
  master: 88c1d2d4197bed815b01a22cadfc6c35aa246dddb553682037a118aebfaa3954
  mizaru_free: 0558216cbab7a9c46f298f4c26e171add9af87d0694988b8a8fe52ee932aa754
  mizaru_plus: cf6f58868c6d9459b3a63bc2bd86165631b3e916bad7f62b578cd9614e0bcb3b
  mizaru_bw: $bw_key
YML

    # 启动进程管理
    procd_open_instance
    procd_set_param command "$prog" -c "$YAML_FILE"
    
    # 5秒内不退则强杀
    procd_set_param term_timeout 5
    
    procd_set_param stderr 1
    procd_set_param stdout 1
    procd_set_param respawn
    procd_close_instance
}

stop_service() {
    # 显式停止时也执行强杀
    killall -9 geph5-client >/dev/null 2>&1
}

service_triggers() {
    procd_add_reload_trigger "$CONF"
}
